# .local/system/81-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/run/lightdm/daniel/xauthority", RUN+="/usr/bin/arandr"
